# live_bot/deploy/hourly_fetch.service
[Unit]
Description=Hourly OI + funding refresh
After=network-online.target

[Service]
Type=oneshot
User=bot
Group=bot
WorkingDirectory=/opt/bot
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/opt/bot
ExecStart=/opt/bot/.venv/bin/python /opt/bot/tools/hourly_fetch.py \
    --symbols BTCUSDT ETHUSDT --intervals 60
