# live_bot/deploy/hourly_fetch.timer
[Unit]
Description=Hourly OI + funding refresh

[Timer]
OnBootSec=3min
OnUnitActiveSec=1h
Persistent=true
RandomizedDelaySec=60

[Install]
WantedBy=timers.target
